# libfreenect-novnc/etc/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :0 -xkb -forever -shared -repeat
autorestart=true